# ops/retriever-nightly/systemd/nrl-harness-batch-hf.timer
# SPDX-FileCopyrightText: Copyright (c) 2026, NVIDIA CORPORATION & AFFILIATES.
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Schedule the NeMo Retriever batch-HF nightly at midnight Eastern

[Timer]
OnCalendar=*-*-* 00:00:00 America/New_York
Persistent=false
AccuracySec=1m
Unit=nrl-harness-batch-hf.service

[Install]
WantedBy=timers.target
