# ops/retriever-nightly/systemd/nrl-harness-batch-hf.service
# SPDX-FileCopyrightText: Copyright (c) 2026, NVIDIA CORPORATION & AFFILIATES.
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Run latest upstream/main NeMo Retriever and ViDoRe v3 nightly session
Documentation=file:/localhome/local-jioffe/NeMo-Retriever-deploy/nightly/ops/retriever-nightly/README.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=local-jioffe
Group=local-jioffe
WorkingDirectory=/localhome/local-jioffe/NeMo-Retriever-deploy/nightly
Environment=HOME=/localhome/local-jioffe
ExecStart=/usr/bin/bash /localhome/local-jioffe/NeMo-Retriever-deploy/nightly/ops/retriever-nightly/run-latest-main.sh
UMask=0077
TimeoutStartSec=infinity
Restart=no
